# crontab
0 0 7 * * /bin/bash /app/kiwix-zim-updater.sh $SCRIPT_FLAGS /zim 2>&1 | tee -a /zim/logs/$(date +%F_%T).log
