# crontab
0 3 1 * * /bin/bash /kiwix-zim-updater.sh $SCRIPT_FLAGS /zim
